# minecraftd-backup.service.in
[Unit]
Description=@GAME@ Server World Backup
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/bin/@INAME@ backup
User=@GAME_USER@
Group=@GAME_USER@

[Install]
WantedBy=multi-user.target
